# edge/deploy/clarus-edge.service
[Unit]
Description=EdgeSentry clarus-edge daemon
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=yohei
WorkingDirectory=/home/yohei/clarus/edge
EnvironmentFile=/home/yohei/clarus/edge/config.env
ExecStart=/usr/local/bin/clarus-edge
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
